# calendar2image/rootfs/etc/services.d/calendar2image/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Take down the calendar2image service
# ==============================================================================

foreground { s6-svc -wD -d -T 5000 /run/service/calendar2image }
