# calendar2image/rootfs/etc/services.d/calendar2image/run
#!/bin/sh
# ==============================================================================
# Start calendar2image service
# ==============================================================================

echo "[INFO] Starting Calendar2Image on port 3000..."

cd /app
exec node src/index.js
