# calendar2image/rootfs/etc/services.d/calendar2image/run
#!/bin/sh
# ==============================================================================
# Start calendar2image service
# ==============================================================================

echo "[INFO] Starting Calendar2Image on port 3000..."

# Export Puppeteer environment variables
export PUPPETEER_SKIP_CHROMIUM_DOWNLOAD=true
export PUPPETEER_EXECUTABLE_PATH=/usr/bin/chromium
export CHROMIUM_PATH=/usr/bin/chromium

# Set config and templates directory to /config
export CONFIG_DIR=/config
export TEMPLATES_DIR=/config/templates

# Get the real path (resolves symlinks and shows actual filesystem path)
export REAL_CONFIG_PATH=$(realpath "${CONFIG_DIR}" 2>/dev/null || echo "${CONFIG_DIR}")

cd /app
exec node src/index.js
