# calendar2image/rootfs/etc/services.d/calendar2image/run
#!/bin/sh
# ==============================================================================
# Start calendar2image service
# ==============================================================================

echo "[INFO] Starting Calendar2Image on port 3000..."

# Export Puppeteer environment variables
export PUPPETEER_SKIP_CHROMIUM_DOWNLOAD=true
export PUPPETEER_EXECUTABLE_PATH=/usr/bin/chromium
export CHROMIUM_PATH=/usr/bin/chromium

# Set config and templates directory to /data (addon_config mount point)
export CONFIG_DIR=/data
export TEMPLATES_DIR=/data/templates

cd /app
exec node src/index.js
